# src/systemd/ubuntu-merges.timer
[Unit]
Description=Ubuntu Merges Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
