# app/mom.service
[Unit]
Description=Run Merge-O-Matic

[Service]
Type=oneshot
User=merge
ExecStart=/home/merge/data/code/cron.daily
TimeoutStartSec=48h
